# mt-logstash-forwarder.conf
[program:mt-logstash-forwarder]
command=/opt/mt-logstash-forwarder/bin/run-container-lsf.sh
stdout_logfile=/var/log/mt-logstash-forwarder.log
stderr_logfile=/var/log/mt-logstash-forwarder.log

